# templates/default/systemd/monit.service.erb
[Unit]
Description=Pro-active monitoring utility for unix systems
After=network.target

[Service]
Type=simple
KillMode=process
ExecStart=<%= node["monit"]["binary_path"] %> -I
ExecStop=<%= node["monit"]["binary_path"] %> quit
ExecReload=<%= node["monit"]["binary_path"] %> reload

[Install]
WantedBy=multi-user.target
